# root/etc/services.d/php-fpm/run
#!/usr/bin/execlineb -P
/usr/local/sbin/php-fpm -y /usr/local/etc/php-fpm.conf -R