# root/etc/services.d/php-fpm/run
#!/usr/bin/with-contenv /usr/bin/execlineb
/usr/local/sbin/php-fpm -y /usr/local/etc/php-fpm.conf -R