# root/etc/services.d/crond/run
#!/usr/bin/env sh
exec gosu root /usr/sbin/crond -fS
